# roles/eth-deploy/templates/eth-beacon.service
[Unit]
Description=eth-beacon

[Service]
Type=simple
ExecStart=/prysm/prysm-beacon-service.sh start
ExecStop=/prysm/prysm-beacon-service.sh stop
Restart=always
RestartSec=2

[Install]
WantedBy=sysinit.target
